# deploy/dev/supervisor-sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
process_name=%(program_name)s
directory=/
autostart=true
autorestart=true
startsecs=10
startretries=3
